# main/storytellerpi-web.service
[Unit]
Description=StorytellerPi Web Interface
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/storytellerpi/main
Environment=PATH=/opt/storytellerpi/venv/bin
ExecStart=/opt/storytellerpi/venv/bin/python /opt/storytellerpi/main/web_interface.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/storytellerpi

# Resource limits
MemoryMax=200M
CPUQuota=50%

[Install]
WantedBy=multi-user.target